# systemd/tft-monitoring.service
[Unit]
Description=TFT Monitoring System (All Services)
After=network.target

[Service]
Type=forking
User=tftuser
Group=tftuser
WorkingDirectory=/opt/tft-monitoring
ExecStart=/opt/tft-monitoring/start_all_corp.sh
ExecStop=/opt/tft-monitoring/stop_all.sh
Restart=on-failure
RestartSec=10
StandardOutput=append:/var/log/tft-monitoring/system.log
StandardError=append:/var/log/tft-monitoring/system.log

# Environment
Environment="PYTHONUNBUFFERED=1"

# Security hardening
PrivateTmp=yes
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/tft-monitoring/logs

# Resource limits
LimitNOFILE=65536
MemoryLimit=8G

[Install]
WantedBy=multi-user.target
