# systemd/tft-metrics-generator.service
[Unit]
Description=TFT Metrics Generator Daemon
After=network.target

[Service]
Type=simple
User=tftuser
Group=tftuser
WorkingDirectory=/opt/tft-monitoring
ExecStart=/usr/bin/python3 /opt/tft-monitoring/metrics_generator_daemon.py
Restart=always
RestartSec=5
StandardOutput=append:/var/log/tft-monitoring/metrics_generator.log
StandardError=append:/var/log/tft-monitoring/metrics_generator.log

# Environment
Environment="PYTHONUNBUFFERED=1"

# Security hardening
PrivateTmp=yes
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes

# Resource limits
LimitNOFILE=4096
MemoryLimit=1G

[Install]
WantedBy=multi-user.target
